# package-openmoko/CONTROL/postrm
#!/bin/sh
update-rc.d $D evopedia remove
