# package-maemo/DEBIAN/postrm
#!/bin/sh

update-rc.d evopedia remove >/dev/null

exit 0
